# templates/atop-restart.timer.epp
[Unit]
Description=Daily restart of atop service

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true

[Install]
WantedBy=timers.target
