# templates/atop-restart.service.epp
[Unit]
Description=Restart atop service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl try-restart atop
